# server/crontab
40 23 * * * bundle exec rake log2sql > /proc/1/fd/1 2>/proc/1/fd/2
50 23 * * * bundle exec rake delete_users --silent > /proc/1/fd/1 2>/proc/1/fd/2
50 23 * * * bundle exec rake log_hours --silent > /proc/1/fd/1 2>/proc/1/fd/2